# cloudlab-client.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/cloudlab-client.sh

[Install]
WantedBy=default.target